# app/config/service.conf
[program:gunicorn]
command=/opt/pypy-env/bin/gunicorn -c gunicorn.ini --log-file=/dev/stdout wsgi
directory=/cross-platform/
;environment=PYPY_GC_MAX_DELTA=200MB
;environment=PATH="/home/ec2-user/cross-platform/.pydev/bin"
redirect_stderr=true
;stdout_logfile=/var/log/wsgi.log
;stdout_logfile_backups = 7
;user=ec2-user
stopsignal=QUIT
autostart=true
autorestart=true
startretries=100
;startsecs=60

;[program:pushman]
;command=/home/ec2-user/gcm-apns/.pyenv/bin/python /home/ec2-user/gcm-apns/gcmapns.py
;directory=/home/ec2-user/gcm-apns/
;redirect_stderr=true
;stdout_logfile=/var/log/pushman.log
